# a6e/resources/etc/local.d/__bye.stop
#!/bin/sh

cd /sys/class/leds/sys && \
echo timer > trigger && echo 888 | tee delay_on > delay_off

sleep 120 && pkill -f "/sbin/openrc shutdown" &
